# deploy/systemd/mobility-warehouse-backup.timer
[Unit]
Description=Run the mobility warehouse verified backup daily

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
RandomizedDelaySec=10m
Unit=mobility-warehouse-backup.service

[Install]
WantedBy=timers.target
